# scripts/com.openclaw.sync-memory-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identification -->
    <key>Label</key>
    <string>com.openclaw.sync-memory-watch</string>
    
    <!-- Program to Run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>--import</string>
        <string>tsx</string>
        <string>/Users/tulioferro/Backend_local Macbook/openclaw_repo/scripts/sync-memory-cloud-bidirectional.ts</string>
        <string>--watch</string>
    </array>
    
    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/tulioferro/Backend_local Macbook/openclaw_repo</string>
    
    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/tulioferro</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    
    <!-- Auto-start and Keep Alive -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    
    <!-- Restart on Failure -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/tulioferro/.openclaw/logs/sync-watch-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tulioferro/.openclaw/logs/sync-watch-stderr.log</string>
    
    <!-- Resource Limits -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
